# etc/s6/logitty-srv/run
#!/bin/execlineb -P
exec agetty -L -8 -n -l /usr/bin/logitty tty2 115200 linux

